# systemd/display_off.service
[Unit]
Description=Run script to turn off monitor.
After=multi-user.target

[Service]
Type=simple
ExecStart=sh /usr/local/bin/turn_display_off
