# systemd/display_on.timer
[Unit]
Description=Schedules turning on monitor.

[Timer]
OnCalendar=7:00

[Install]
WantedBy=timers.target
