# pic_frame/systemd/display_off.timer
[Unit]
Description=Schedules turning off monitor.

[Timer]
OnCalendar=22:00

[Install]
WantedBy=timers.target
